# airmon.service
[Unit]
Description=Airmon Service

[Service]
Type=simple
ExecStart=/usr/bin/airmon
Restart=always

[Install]
WantedBy=multi-user.target
